# infra/systemd/docker-dashboard.service
[Unit]
Description=Home Dashboard Docker Container
After=docker.service
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/home/pi/home-dashboard
ExecStart=/usr/bin/docker-compose up
ExecStop=/usr/bin/docker-compose down
Restart=on-failure
RestartSec=10
User=pi

[Install]
WantedBy=multi-user.target
