# infra/systemd/docker-dashboard.service
[Unit]
Description=Home Dashboard Docker Container
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/home-dashboard/docker
ExecStart=/usr/bin/docker compose -f %h/home-dashboard/docker/docker-compose.yml up
ExecStop=/usr/bin/docker compose -f %h/home-dashboard/docker/docker-compose.yml down
Restart=on-failure
RestartSec=10
User=pi
EnvironmentFile=-%h/home-dashboard/.env

[Install]
WantedBy=multi-user.target
